# supercronic/crontab
# https://github.com/aptible/supercronic/tree/master/cronexpr
0 0 */3 * * * * bash -c "python -u /opt/app/sync.py"
