# supercronic/crontab
# https://github.com/aptible/supercronic/tree/master/cronexpr
* * * * * * * bash -c "python -u /opt/app/sync.py"
